# setup/roles/app/templates/service.target.j2
[Unit]
Description={{ app_user }} Workers
Requires={{ app_user }}@{{ port }}.service {{ app_user }}@{{ port_2 }}.service 

[Install]
WantedBy=multi-user.target