# setup/roles/app/templates/service.target.j2
[Unit]
Description={{ app_user }} Workers
Requires={{ app_user }}@{{ app_port }}.service

[Install]
WantedBy=multi-user.target
